# com.revenuecat.rev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.revenuecat.rev</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ivanma/Desktop/agents/RevCat/run_agent.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ivanma/Desktop/agents/RevCat</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ivanma/Desktop/agents/RevCat/logs/rev.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ivanma/Desktop/agents/RevCat/logs/rev.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
